# roles/celerycam/templates/supervisor_celerycam.conf.j2
[program:celerycam]
command={{ celerycam_scripts_dir }}/celerycam_start

autostart=true
autorestart=true

user={{ celerycam_user }}

stdout_logfile={{ celerycam_log_file }}
redirect_stderr = true