# roles/elasticsearch/templates/supervisor_elasticsearch.conf.j2
[program:elasticsearch]
command=elasticsearch
autostart=true
autorestart=true
